# ansible/templates/docker.service.j2
[Service]
LimitNOFILE=65536
LimitNPROC=32768
LimitCORE=0
LimitMEMLOCK=infinity 